# ansible/roles/jeffreydd.minisub.vessel/templates/systemd.service.j2
[Unit]
Description={{ svc_name }}

[Service]
ExecStart={{ svc_cmd }}

[Install]
WantedBy=default.target